# system_files/usr/lib/systemd/system/fail2ban.service.d/override.conf
# By default, Fail2ban crashes due to not being able
# to write the SQLite database.
[Service]
ExecStart=
ExecStart=/usr/bin/fail2ban-server -xf -vvvv --logtarget=sysout start
